# rootfs/etc/s6-overlay/s6-rc.d/check-plugin-updates/run
#!/bin/sh
exec /etc/s6-overlay/scripts/check-plugin-updates
